# cron/crontab
* * * * * /usr/scripts/restartmodem.sh >> /var/log/restartmodem.log 2>&1
